# ansible/roles/provision_host/templates/podman-prune.timer.j2
[Unit]
Description=Daily rootless container storage prune

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=20min
Persistent=true

[Install]
WantedBy=timers.target
