# ansible/roles/provision_host/templates/runner-health-watch.timer.j2
[Unit]
Description=Poll runner cgroup health every {{ health_watch_interval }}

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ health_watch_interval }}
AccuracySec=15s

[Install]
WantedBy=timers.target
